# rc.local
#!/bin/sh -e
sslocal -c /etc/shadowsocks.json -d start
./etc/client_linux_amd64 -c /etc/client-config.json >/dev/null 2>&1

exit 0

